# systemd/.config/systemd/user/bpeak.target
[Unit]
Description=Stuff for bpeak
After=xorg.target

[Install]
Alias=default.target
