# systemd/.config/systemd/user/me.target
[Unit]
Description=Stuff for shivalva
After=xorg.target

[Install]
Alias=default.target
